# clockon.timer
[Unit]
Description=MicroKeeper Clock on Automator Timer

[Timer]
OnCalendar=Mon..Fri 9:00

[Install]
WantedBy=timers.target
