# clockoff.timer
[Unit]
Description=MicroKeeper Clock off Automator Timer

[Timer]
OnCalendar=Mon..Fri 17:00

[Install]
WantedBy=timers.target
