# clockon.service
[Unit]
Description=MicroKeeper Clock on Automator Service
Wants=clockon.timer

[Service]
Type=oneshot
ExecStart=%h/bin/clockon.py
