# clockon.timer
[Unit]
Description=MicroKeeper Clock on Automator Timer
Requires=clockon.service

[Timer]
OnCalendar=Mon..Fri 9:00

[Install]
WantedBy=timers.target
